# gnome-extra/iio-sensor-proxy/files/watch-sensors.service
[Unit]
Description=Poll Sensors Service

[Service]
Type=exec
ExecStart=monitor-sensor

[Install]
WantedBy=multi-user.target
